# systemd/.config/systemd/user/battery-warning.service
[Unit]
Description=Desktop Warning of Low Battery
[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery-warning
